# ops/supervisor.conf
[program:noelwilson2018]
directory=/opt/noelwilson_2018/current/
command=/opt/noelwilson_2018/current/
user=root
redirect_stderr=true
stdout_logfile=/var/log/noelwilson2018.log
