# ops/supervisor.conf
[program:noelwilson]
directory=/opt/app/current/
command=make run-be
user=root
redirect_stderr=true
stdout_logfile=/var/log/noelwilson.log
